# backend/aria-backend.service
[Unit]
Description=Aria Backend
Documentation=https://github.com/chaitea321/MusicAppIOS
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=eugene
WorkingDirectory=/home/eugene/MusicAppIOS/backend
Environment=MAX_CACHE_GB=5
Environment=CACHE_EVICT_GRACE_SECONDS=300
# Free the ports if a previous instance didn't exit cleanly, then start.
ExecStartPre=/bin/sh -c "for p in $(fuser 8000/tcp 8443/tcp 2>/dev/null); do kill -9 $p; done; sleep 1"
ExecStart=/home/eugene/MusicAppIOS/backend/.venv/bin/python3 -m uvicorn app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
